# tests/tests/tier0/bluechi-restart-unit/systemd/simple.service
[Unit]
Description=Run a service for 15 seconds

[Service]
Type=simple
ExecStart=/bin/true
RemainAfterExit=yes
