# tests/tests/tier0/proxy-service-stop-bluechi-dep/systemd/requesting.service
[Unit]
Wants=bluechi-proxy@node-bar_simple.service
After=bluechi-proxy@node-bar_simple.service

[Service]
ExecStart=/bin/true
RemainAfterExit=yes
