# tests/tests/tier0/bluechi-reload-unit/systemd/service-for-reload.service
[Unit]
Description=A very simple service
[Service]
Type=simple
ExecStart=/bin/true
RemainAfterExit=yes
ExecReload=systemctl start simple.service