# tests/tests/tier0/bluechi-reload-unit/systemd/simple.service
[Unit]
Description=Service For Integration Test

[Service]
Type=simple
ExecStart=/bin/true
RemainAfterExit=yes
