# tests/tests/tier0/proxy-service-stop-requesting-with-unneeded/systemd/simple.service
[Unit]
Description=Just sleeping
StopWhenUnneeded=yes

[Service]
Type=simple
ExecStart=/bin/true
RemainAfterExit=yes
