# tests/tests/tier0/proxy-service-fails-on-execstart/systemd/simple.service
[Unit]
Description=Attempting to just sleep, but fail

[Service]
Type=simple
ExecStart=exit 1
RemainAfterExit=yes
