# tests/tests/tier0/bluechi-monitor-service/systemd/simple.service
[Unit]
Description=Run a simple service

[Service]
Type=simple
ExecStart=/bin/true
RemainAfterExit=yes
