# tests/tests/tier0/bluechi-freeze-and-thaw-service/systemd/simple.service
[Unit]
Description=Just sleeping

[Service]
Type=simple
ExecStart=/bin/sleep 10
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
