# tests/tests/tier0/monitor-listener-removed-as-peer/systemd/simple.service
[Unit]
Description=Just being true once

[Service]
Type=simple
ExecStart=/bin/true
