# tests/tests/tier0/bluechi-daemon-reload/systemd/simple.service
[Service]
ExecStart=/s
RemainAfterExit=yes