# tests/tests/tier0/bluechi-enable-disable-service/systemd/simple.service
[Unit]
Description=I am truth

[Service]
Type=simple
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
