# tests/tests/tier0/proxy-service-start/systemd/requesting.service
[Unit]
Wants=hirte-proxy@node-bar_simple.service
After=hirte-proxy@node-bar_simple.service

[Service]
ExecStart=/bin/true
RemainAfterExit=yes
