# tests/tests/tier0/proxy-service-start/systemd/simple.service
[Unit]
Description=Just sleeping

[Service]
Type=simple
ExecStart=/bin/true
RemainAfterExit=yes
