# tests/tests/tier0/proxy-service-fails-on-typo-in-file/systemd/simple.service
[Unit]
Description=Just sleeping

[Sevice]
Type=simple
ExecStart=/bin/true
RemainAfterExit=yes
