# tests/tests/tier0/bluechi-job-cancel/systemd/simple.service
[Unit]
Description=Just sleeping

[Service]
Type=simple
ExecStartPre=/bin/sleep 20
ExecStart=/bin/true
RemainAfterExit=yes
