# tests/tests/tier0/monitor-system-status/systemd/monitor.service
[Unit]
Description=Monitor BlueChi system

[Service]
Type=simple
ExecStart=python3 /tmp/system-monitor.py

[Install]
WantedBy=multi-user.target
